# deploy/systemd/exchange-rq-gas.service
[Unit]
Description=Exchange RQ worker - gas
After=network-online.target redis.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/exchange-web/backend
EnvironmentFile=/opt/exchange-web/backend/.env
User=exchange
Group=exchange
ExecStart=/opt/exchange-web/.venv/bin/python scripts/start_rq_worker.py gas
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
