# deploy/systemd/exchange-api.service
[Unit]
Description=Exchange FastAPI API
After=network-online.target redis.service
Wants=network-online.target redis.service
PartOf=exchange-backend.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/opt/exchange-web/backend
EnvironmentFile=/opt/exchange-web/backend/.env
Environment=PYTHONUNBUFFERED=1
Environment=EMBED_BACKGROUND_LOOPS_IN_API=0
Environment=ENABLE_SPOT_AUTO_MATCH_IN_API=0
Environment=ENABLE_CONTRACT_LIMIT_ORDER_JOB=0
User=exchange
Group=exchange
ExecStart=/opt/exchange-web/.venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000 --access-log
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
KillMode=mixed
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
UMask=0027

[Install]
WantedBy=multi-user.target
